# archiso/airootfs/etc/skel/.config/autostart/umaos-welcome.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=UmaOS Welcome (First Login)
Comment=Launch UmaOS welcome wizard on first login
Exec=/usr/local/bin/umao-welcome
Terminal=false
OnlyShowIn=KDE;
X-KDE-autostart-phase=2
