# archiso/airootfs/usr/lib/systemd/user/umao-wallpaper.service
[Unit]
Description=UmaOS Dynamic Wallpaper Switcher

[Service]
Type=oneshot
ExecStart=/usr/local/bin/umao-wallpaper-timer
