# archiso/airootfs/etc/profile.d/umaos-welcome.sh
#!/usr/bin/env bash
if [[ -n "${PS1-}" ]]; then
  printf '\nUmaOS live session ready. Calamares auto-launches in KDE; run `umao-install` anytime.\n'
  printf 'Install Uma Musume from desktop via `Install Uma Musume.sh`.\n\n'
fi
